# systemd/glances-systemd/glancesd.service
[Unit]
Description=Glances daemon
After=network.target

[Service]
ExecStart=/usr/bin/glances -s

[Install]
WantedBy=multi-user.target
